# contrib/init/syscoind.conf
description "Syscoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SYSCOIND_BIN="/usr/bin/syscoind"
env SYSCOIND_USER="syscoin"
env SYSCOIND_GROUP="syscoin"
env SYSCOIND_PIDDIR="/var/run/syscoind"
# upstart can't handle variables constructed with other variables
env SYSCOIND_PIDFILE="/var/run/syscoind/syscoind.pid"
env SYSCOIND_CONFIGFILE="/etc/syscoin/syscoin.conf"
env SYSCOIND_DATADIR="/var/lib/syscoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # syscoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SYSCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run syscoind."
        echo "The setting must appear in $SYSCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Syscoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SYSCOIND_PIDDIR"
    chmod 0755 "$SYSCOIND_PIDDIR"
    chown $SYSCOIND_USER:$SYSCOIND_GROUP "$SYSCOIND_PIDDIR"
    chown $SYSCOIND_USER:$SYSCOIND_GROUP "$SYSCOIND_CONFIGFILE"
    chmod 0660 "$SYSCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SYSCOIND_PIDFILE" \
    --chuid $SYSCOIND_USER:$SYSCOIND_GROUP \
    --exec "$SYSCOIND_BIN" \
    -- \
    -pid="$SYSCOIND_PIDFILE" \
    -conf="$SYSCOIND_CONFIGFILE" \
    -datadir="$SYSCOIND_DATADIR" \
    -disablewallet \
    -daemon

