# contrib/init/talkcoind.conf
description "Talkcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env TALKCOIND_BIN="/usr/bin/talkcoind"
env TALKCOIND_USER="talkcoin"
env TALKCOIND_GROUP="talkcoin"
env TALKCOIND_PIDDIR="/var/run/talkcoind"
# upstart can't handle variables constructed with other variables
env TALKCOIND_PIDFILE="/var/run/talkcoind/talkcoind.pid"
env TALKCOIND_CONFIGFILE="/etc/talkcoin/talkcoin.conf"
env TALKCOIND_DATADIR="/var/lib/talkcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # talkcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$TALKCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run talkcoind."
        echo "The setting must appear in $TALKCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Talkcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$TALKCOIND_PIDDIR"
    chmod 0755 "$TALKCOIND_PIDDIR"
    chown $TALKCOIND_USER:$TALKCOIND_GROUP "$TALKCOIND_PIDDIR"
    chown $TALKCOIND_USER:$TALKCOIND_GROUP "$TALKCOIND_CONFIGFILE"
    chmod 0660 "$TALKCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$TALKCOIND_PIDFILE" \
    --chuid $TALKCOIND_USER:$TALKCOIND_GROUP \
    --exec "$TALKCOIND_BIN" \
    -- \
    -pid="$TALKCOIND_PIDFILE" \
    -conf="$TALKCOIND_CONFIGFILE" \
    -datadir="$TALKCOIND_DATADIR" \
    -disablewallet \
    -daemon

